# timefree/crontab
SHELL=/bin/bash
PATH=/usr/local/bin:/bin:/usr/bin
*/10  *  * * * bundle exec ruby ./sqs_polling.rb
