# timefree/crontab
SHELL=/bin/bash
PATH=/usr/local/bin:/bin:/usr/bin
*/10  *  * * * cd /work && make sqs
